# src/fty-outage.service.in
[Unit]
Description=fty-outage service
After=malamute.service network.target
Requires=malamute.service network.target
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/fty
EnvironmentFile=-/etc/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-outage @sysconfdir@/@PACKAGE@/fty-outage.cfg
Restart=always

[Install]
WantedBy=bios.target
